# tradeTip.cron
50 14 * * 1-5 daemon sh /data/tradeTip/run.sh &> /tmp/tradeTipCron.log
