# tradeTip.cron
50 14 * * 1-5 root sh /data/tradeTip/run.sh &> /dev/null
